# SimpleSeer/paster_templates/seer_project/etc/supervisor_mongodb.conf
[program:mongodb]
command=/usr/local/bin/mongod --config /etc/mongodb.conf
process_name=%(program_name)s
autostart=True
priority=1
startsecs=5
startretries=20
