# SimpleSeer/paster_templates/seer_project/settings/default/conf.d/supervisor_mongodb.conf
[program:mongodb]
command=/usr/local/bin/mongod --config /etc/mongodb.conf
process_name=%(program_name)s
autostart=True
priority=2
startsecs=10
startretries=20
